# systemd/vps-safe-cleanup.timer
[Unit]
Description=Run VPS safe cleanup hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
AccuracySec=5min
Persistent=true
Unit=vps-safe-cleanup.service

[Install]
WantedBy=timers.target
